# .oblivious/services/server/run
#!/bin/sh
python /code/uvicorn_serve.py
